# roles/edpm_instanceha_monitoring/templates/instanceha-heartbeat.service.j2
# {{ ansible_managed }}
[Unit]
Description=InstanceHA heartbeat sender
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/instanceha-heartbeat.py
